# service/hm-mb-control.timer
[Unit]
Description=Run hm-mb-control after 1 minute delay

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target